# IPCAM_Reference/Application/mozart3s_Kilrogg_r47602_IPCam/apps/initscripts/init.d-synctime
#!/bin/sh
# /etc/init.d/synctime: sync system time to RTC
# chkconfig: 2345 08 60

PATH=/sbin:/bin:/usr/bin

RETVAL=0

prog="synctime"
args="-s -u"

start(){
    echo -n "Starting $prog: "
	hwclock $args
	echo
}

stop(){
    echo -n "Stopping $prog: "
	echo "."
}

# See how we were called.
case "$1" in
    start)
	start
	;;
    stop)
	stop
	;;
    restart)
    stop
    start
	;;
    *)
	echo "Usage: $0 {start|stop|restart}"
	exit 1
esac

exit 0
