# ASC88xx_SDK/LSP/mozart3_rootfs-1.18/static/etc/inittab
console::sysinit:-/etc/init.d/rcS
::respawn:/sbin/getty -L ttyS 38400 vt100
::restart:/sbin/init
